# actors/profiles/fin8.yaml
id: fin8
name: FIN8
mitre_id: G0061
aliases:
    - Syssphinx
type: criminal
description: '[FIN8](https://attack.mitre.org/groups/G0061) is a financially motivated threat group that has been active since at least January 2016, and known for targeting organizations in the hospitality, retail, entertainment, insurance, technology, chemical, and financial sectors. In June 2021, security researchers detected [FIN8](https://attack.mitre.org/groups/G0061) switching from targeting point-of-sale (POS) devices to distributing a number of ransomware variants.(Citation: FireEye Obfuscation June ...'
ttps:
    - id: T1078
      name: Valid Accounts
    - id: T1048.003
      name: Exfiltration Over Unencrypted Non-C2 Protocol
    - id: T1033
      name: System Owner/User Discovery
    - id: T1518.001
      name: Security Software Discovery
    - id: T1021.001
      name: Remote Desktop Protocol
    - id: T1003.001
      name: LSASS Memory
    - id: T1588.002
      name: Tool
    - id: T1204.002
      name: Malicious File
    - id: T1588.003
      name: Code Signing Certificates
    - id: T1068
      name: Exploitation for Privilege Escalation
    - id: T1546.003
      name: Windows Management Instrumentation Event Subscription
    - id: T1566.002
      name: Spearphishing Link
    - id: T1053.005
      name: Scheduled Task
    - id: T1204.001
      name: Malicious Link
    - id: T1102
      name: Web Service
    - id: T1027.010
      name: Command Obfuscation
    - id: T1070.004
      name: File Deletion
    - id: T1566.001
      name: Spearphishing Attachment
    - id: T1071.001
      name: Web Protocols
    - id: T1021.002
      name: SMB/Windows Admin Shares
    - id: T1685.005
      name: Clear Windows Event Logs
    - id: T1560.001
      name: Archive via Utility
    - id: T1074.002
      name: Remote Data Staging
    - id: T1105
      name: Ingress Tool Transfer
    - id: T1082
      name: System Information Discovery
    - id: T1059.001
      name: PowerShell
    - id: T1059.003
      name: Windows Command Shell
    - id: T1573.002
      name: Asymmetric Cryptography
    - id: T1055.004
      name: Asynchronous Procedure Call
    - id: T1018
      name: Remote System Discovery
    - id: T1486
      name: Data Encrypted for Impact
    - id: T1482
      name: Domain Trust Discovery
    - id: T1112
      name: Modify Registry
    - id: T1134.001
      name: Token Impersonation/Theft
    - id: T1016.001
      name: Internet Connection Discovery
    - id: T1047
      name: Windows Management Instrumentation
software:
    - Ping
    - BADHATCH
    - PUNCHBUGGY
    - Ragnar Locker
    - PUNCHTRACK
    - dsquery
    - Net
    - Nltest
    - Sardonic
    - PsExec
    - Impacket
confidence: high
