# actors/profiles/storm-0501.yaml
id: storm-0501
name: Storm-0501
mitre_id: G1053
type: criminal
description: '[Storm-0501](https://attack.mitre.org/groups/G1053) is a financially motivated cyber criminal group that uses commodity and open-source tools to conduct ransomware operations. [Storm-0501](https://attack.mitre.org/groups/G1053) has been active since 2021 and has previously been affiliated with Sabbath Ransomware and other Ransomware-as-a-Service (RaaS) variants such as Hive, [BlackCat](https://attack.mitre.org/software/S1068), Hunters International, [LockBit 3.0](https://attack.mitre.org/softwar...'
ttps:
    - id: T1219.002
      name: Remote Desktop Software
    - id: T1537
      name: Transfer Data to Cloud Account
    - id: T1490
      name: Inhibit System Recovery
    - id: T1484.001
      name: Group Policy Modification
    - id: T1530
      name: Data from Cloud Storage
    - id: T1059.001
      name: PowerShell
    - id: T1485
      name: Data Destruction
    - id: T1053.005
      name: Scheduled Task
    - id: T1087.004
      name: Cloud Account
    - id: T1003.006
      name: DCSync
    - id: T1526
      name: Cloud Service Discovery
    - id: T1567.002
      name: Exfiltration to Cloud Storage
    - id: T1059.009
      name: Cloud API
    - id: T1021.007
      name: Cloud Services
    - id: T1021.006
      name: Windows Remote Management
    - id: T1190
      name: Exploit Public-Facing Application
    - id: T1057
      name: Process Discovery
    - id: T1518.001
      name: Security Software Discovery
    - id: T1486
      name: Data Encrypted for Impact
    - id: T1657
      name: Financial Theft
    - id: T1078.004
      name: Cloud Accounts
    - id: T1218.010
      name: Regsvr32
    - id: T1555.006
      name: Cloud Secrets Management Stores
    - id: T1555.005
      name: Password Managers
    - id: T1484.002
      name: Trust Modification
    - id: T1580
      name: Cloud Infrastructure Discovery
    - id: T1482
      name: Domain Trust Discovery
    - id: T1578.003
      name: Delete Cloud Instance
    - id: T1082
      name: System Information Discovery
    - id: T1027.002
      name: Software Packing
    - id: T1614.001
      name: System Language Discovery
    - id: T1552.004
      name: Private Keys
    - id: T1036.004
      name: Masquerade Task or Service
    - id: T1098.003
      name: Additional Cloud Roles
    - id: T1087.002
      name: Domain Account
    - id: T1218.011
      name: Rundll32
    - id: T1587.003
      name: Digital Certificates
    - id: T1588.006
      name: Vulnerabilities
    - id: T1556.009
      name: Conditional Access Policies
    - id: T1003
      name: OS Credential Dumping
    - id: T1110
      name: Brute Force
    - id: T1098.001
      name: Additional Cloud Credentials
software:
    - Impacket
    - Tasklist
    - Cobalt Strike
    - Embargo
    - Rclone
    - Nltest
    - Net
    - AADInternals
confidence: high
